# base/rootfs/etc/services.d/confd/run
#!/usr/bin/with-contenv bash
set -vxe

if [ -n "${CONFD_NOOP}" ] && [ "false" != "${CONFD_NOOP}" ]
then
  confd-render-templates.sh
else
  confd-render-templates.sh --continuous
fi
